# cron/jobs.cron
0 * * * * mailsync
0 0 * * * make -C ~/dotfiles update
