# cron/.config/cron/crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/Users/shreeram/.nix-profile/bin/

*/15 * * * * $HOME/.local/bin/mail-autosync
*/15 * * * * $HOME/.nix-profile/bin/calcurse-caldav
