# cron/.config/cron/crontab
*/15 * * * * $HOME/.local/bin/mail-autosync
*/15 * * * * $HOME/.nix-profile/bin/calcurse-caldav
