# cron/.config/cron/crontab
*/15 * * * * $HOME/.local/bin/mail-autosync
